# packages/client/keenyspace/templates/launchd_com.keenyspace.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.keenyspace.daemon</string>
  <key>ProgramArguments</key>
  <array>
    <string>__KEENYSPACE_BIN__</string>
    <string>daemon</string>
    <string>start</string>
    <string>--foreground</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>__HOME__/.local/state/keenyspace/daemon.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/.local/state/keenyspace/daemon.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>XDG_CONFIG_HOME</key>
    <string>__HOME__/.config</string>
    <key>XDG_STATE_HOME</key>
    <string>__HOME__/.local/state</string>
  </dict>
</dict>
</plist>
