# examples/systemd/keenyspace.service
[Unit]
Description=KeenySpace background daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=__KEENYSPACE_BIN__ daemon start --foreground
Restart=on-failure
RestartSec=5
Environment=XDG_CONFIG_HOME=%h/.config
Environment=XDG_STATE_HOME=%h/.local/state
StandardOutput=append:%h/.local/state/keenyspace/daemon.log
StandardError=append:%h/.local/state/keenyspace/daemon.log

[Install]
WantedBy=default.target
